# systemd/export-pkgs.service
[Unit]
Description=Export Arch package list to dotfiles

[Service]
Type=oneshot
ExecStart=/home/estebanpuyanas/dotfiles/export-packages.sh

